# ntpd.service
[Unit]
Description=Network Time Protocol Deamon

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStartPre=
ExecStart=/usr/bin/ntpd
ExecStartPost=
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target
