# ttyd.conf
[program:ttyd]
command=/usr/bin/ttyd -p 443 -c gt:a! bash
autorestart=true
startretries=99999999
